# systemd/disk-sentinel.service
[Unit]
Description=Daily disk sentinel check

[Service]
Type=oneshot
ExecStart=/home/raizen/.local/bin/disk-sentinel.sh

[Install]
WantedBy=multi-user.target
